# util/epiphyte-aur.timer
[Unit]
Description=aur pkgbuilds

[Timer]
OnCalendar=daily
Unit=epiphyte-aur.service

[Install]
WantedBy=timers.target
