# util/epiphyte-autobuild.service
[Unit]
Description=autobuild epiphyte pkgbuilds

[Service]
Type=oneshot
User=autobuild
Group=autobuild
ExecStart=/usr/bin/epiphyte-autobuild
