# util/epiphyte-autobuild.service
[Unit]
Description=autobuild epiphyte pkgbuilds

[Service]
Type=oneshot
ExecStart=/usr/bin/epiphyte-autobuild
